# deb/debian/postrm
#!/bin/sh
# postrm script for vph-secproxy
#
# see: dh_installdeb(1)

#set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <overwriter>
#          <overwriter-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

remove_from_autostart() {
  update-rc.d -f vphsecproxy remove
  update-rc.d -f vphsecagent remove

  exit 0
}

remove_libs() {
  rm -rf /usr/share/vph-secproxy

  exit 0
}

remove_configs() {
  rm -rf /etc/vph-secproxy

  exit 0
}

case "$1" in
    remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)      
      remove_from_autostart
      remove_libs
    ;;

    purge)
      remove_configs
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
